# pwr_test/power_cycle_check.service
[Unit]
Description=Power Cycle Check Service
After=network.target

[Service]
ExecStart=/usr/local/bin/power_cycle_check.sh
Restart=no

[Install]
WantedBy=multi-user.target
